# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
pidfile=/var/run/supervisord/supervisord.pid
loglevel=info
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:general]
command=/app/meta-agent/pi/pi-test.sh --no-context-files --name general --provider minimax --model MiniMax-M3 \
  -e /app/meta-agent/meta-agent-config/extensions/planning-with-files \
  -e /app/meta-agent/meta-agent-config/v1/mission-control/extensions/mission-control/index.ts \
  -e /app/meta-agent/meta-agent-config/v1/sub-agent/src/src/extension/index.ts \
  -e /app/meta-agent/meta-agent-config/v1/sub-agent-context/extensions/sac/index.ts \
  -e /app/meta-agent/meta-agent-config/v1/communication/extensions/communication/index.ts \
  -e /app/meta-agent/meta-agent-config/v1/superhive/extensions/superhive/index.ts \
  -e /app/meta-agent/meta-agent-config/v1/integrations/mc-sac.ts \
  -e /app/meta-agent/meta-agent-config/v1/integrations/planning-mc.ts \
  -e /app/meta-agent/meta-agent-config/v1/integrations/sac-subagent.ts \
  -e /app/meta-agent/meta-agent-config/v1/integrations/comm-perm.ts \
  -e /app/meta-agent/meta-agent-config/v1/integrations/comm-subagent.ts \
  -e /app/meta-agent/meta-agent-config/v1/integrations/comm-planning.ts \
  --skill /app/meta-agent/meta-agent-config/skills/constitutions/00-CONSTITUTION-general.md \
  --skill /app/meta-agent/meta-agent-config/skills/personas/10-PERSONA-general.md \
  --skill /app/meta-agent/meta-agent-config/v1/identity/SKILL.md \
  --skill /app/meta-agent/meta-agent-config/v1/docs/SKILL.md \
  --skill /app/meta-agent/meta-agent-config/v1/planning/SKILL.md \
  --skill /app/meta-agent/meta-agent-config/v1/browser/SKILL.md \
  --skill /app/meta-agent/meta-agent-config/v1/mission-control/SKILL.md \
  --skill /app/meta-agent/meta-agent-config/v1/sub-agent/SKILL.md \
  --skill /app/meta-agent/meta-agent-config/v1/sub-agent-context/SKILL.md \
  --skill /app/meta-agent/meta-agent-config/skills/planning-with-files.md \
  --skill /app/meta-agent/meta-agent-config/skills/browser.md \
  --append-system-prompt /app/meta-agent/meta-agent-config/prompts/extra-instructions.md
directory=/root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisord/general-stdout.log
stderr_logfile=/var/log/supervisord/general-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
user=root
environment=HOME="/root",SUPERHIVE_API_KEY="%(ENV_SUPERHIVE_API_KEY)s",PI_CODING_AGENT_DIR="/app/meta-agent/meta-agent-config",PI_CODING_AGENT_SESSION_DIR="/root/.pi/agent/sessions",PI_SKIP_VERSION_CHECK="1",NODE_ENV="production",MINIMAX_API_KEY="%(ENV_MINIMAX_API_KEY)s",ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s",GEMINI_API_KEY="%(ENV_GEMINI_API_KEY)s",OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",DEEPSEEK_API_KEY="%(ENV_DEEPSEEK_API_KEY)s",GROQ_API_KEY="%(ENV_GROQ_API_KEY)s",MISTRAL_API_KEY="%(ENV_MISTRAL_API_KEY)s",OPENROUTER_API_KEY="%(ENV_OPENROUTER_API_KEY)s",TOGETHER_API_KEY="%(ENV_TOGETHER_API_KEY)s",PI_OFFLINE="%(ENV_PI_OFFLINE)s"
priority=200
