# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log
pidfile=/var/run/supervisord/supervisord.pid
loglevel=info
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisor，一只]
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info

[program:superhive]
command=node /app/v1/superhive/dist/host.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisord/superhive-stdout.log
stderr_logfile=/var/log/supervisord/superhive-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
user=root
environment=HOME="/root",SUPERHIVE_API_KEY="%(ENV_SUPERHIVE_API_KEY)s",NODE_ENV="production"
priority=100
